# packaging/kindle/ferrink.conf
# shellcheck shell=sh disable=SC2093
start on langpicker_ready
stop on stopping x or ota-update

script
    exec /usr/bin/env -i \
        PATH=/usr/sbin:/usr/bin:/sbin:/bin \
        /var/local/ferrink/boot
end script

post-stop script
    PATH=/usr/sbin:/usr/bin:/sbin:/bin

    x_status=$(/sbin/initctl status x 2>/dev/null || true)
    ota_status=$(/sbin/initctl status ota-update 2>/dev/null || true)
    case "$x_status:$ota_status" in
        "x start/running"*:"ota-update stop/waiting"*) ;;
        *) exit 0 ;;
    esac

    gui_status=$(/sbin/initctl status lab126_gui 2>/dev/null || true)
    case "$gui_status" in
        "lab126_gui start/running"*) exit 0 ;;
    esac

    /usr/bin/logger -t ferrink-boot -- "starting stock fallback" 2>/dev/null || true
    /sbin/initctl start lab126_gui FERRINK_FALLBACK=1 || true
end script
